# system_files/usr/share/xdg/autostart/hyprland-portal.desktop
[Desktop Entry]
Type=Application
Name=Hyprland Desktop Portal
Exec=/usr/libexec/xdg-desktop-portal-hyprland

